# supervisor.conf
[supervisord]
nodaemon=true

[program:registry]
command=docker-registry

[program:nginx]
command=/usr/local/bin/run-nginx
